# deploy/flowin-queue-worker.conf
; Supervisor program that processes the queue used to deliver leads to Flowin.
; Copy to /etc/supervisor/conf.d/flowin-queue-worker.conf on the VPS, then:
;   sudo supervisorctl reread && sudo supervisorctl update && sudo supervisorctl start flowin-queue-worker:*
;
; Adjust `command` path, `user`, and `numprocs` to match the server.

[program:flowin-queue-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/laravel/artisan queue:work --queue=default --tries=5 --max-time=3600 --sleep=3
directory=/var/www/laravel
autostart=true
autorestart=true
stopwaitsecs=3600
startsecs=0
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/laravel/storage/logs/flowin-queue-worker.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
